# heatpump-bridge/deploy/heatpump-bridge.service
# @purpose: systemd unit for the bridge on the Pi 5. Install:
#   sudo cp deploy/heatpump-bridge.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now heatpump-bridge
[Unit]
Description=Heat pump Modbus bridge (A2W Control)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/a2w-control/heatpump-bridge
Environment=BRIDGE_CONFIG=/home/pi/a2w-control/heatpump-bridge/config.yaml
# 0.0.0.0 = reachable from the home LAN (http://heatpump-pi.local:8000).
# Once Cloudflare Tunnel is set up, change to --host 127.0.0.1 so the tunnel
# (with its email-OTP auth) is the only way in.
ExecStart=/home/pi/a2w-control/heatpump-bridge/.venv/bin/uvicorn bridge.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
